# pkg/postinstall
#!/bin/bash
# vim: tw=100:

# This file is part of Background Music.
#
# Background Music is free software: you can redistribute it and/or
# modify it under the terms of the GNU General Public License as
# published by the Free Software Foundation, either version 2 of the
# License, or (at your option) any later version.
#
# Background Music is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with Background Music. If not, see <http://www.gnu.org/licenses/>.

#
# postinstall
#
# Copyright © 2017-2022, 2024 Kyle Neideck
# Copyright © 2023 modue sp. z o.o.
#

# Make sure we use the built-in versions of programs, for consistency. Probably not necessary
# because the installer will set $PATH in a similar way, but it doesn't hurt.
PATH=/bin:/sbin:/usr/bin:/usr/sbin; export PATH

function log {
    logger "$@"
    echo "$@"
}

function log_output {
    while read line
    do
        log "$line"
    done < /dev/stdin
}

capture_command_with_timeout() {
    local timeout_seconds="$1"
    local output_var_name="$2"
    shift 2

    local had_errexit=0
    [[ $- == *e* ]] && had_errexit=1
    set +e

    local tmpfile
    tmpfile="$(mktemp -t bgm-timeout.XXXXXX)"
    "$@" >"${tmpfile}" 2>&1 &
    local command_pid=$!
    local exit_status=0
    local timed_out=0
    local waited=0

    while kill -0 "${command_pid}" > /dev/null 2>&1; do
        if [[ ${waited} -ge ${timeout_seconds} ]]; then
            timed_out=1
            kill "${command_pid}" > /dev/null 2>&1 || true
            sleep 1
            kill -9 "${command_pid}" > /dev/null 2>&1 || true
            break
        fi

        sleep 1
        waited=$((waited + 1))
    done

    if [[ ${timed_out} -eq 0 ]]; then
        wait "${command_pid}"
        exit_status=$?
    else
        wait "${command_pid}" > /dev/null 2>&1 || true
        exit_status=124
    fi

    local captured_output
    captured_output="$(cat "${tmpfile}")"
    rm -f "${tmpfile}"
    printf -v "${output_var_name}" '%s' "${captured_output}"

    if [[ ${had_errexit} -eq 1 ]]; then
        set -e
    fi

    return "${exit_status}"
}

verify_installed_bundle_for_hal_restart() {
    local bundle_path="$1"
    local bundle_label="$2"
    local verify_output=""
    local details_output=""

    if ! verify_output="$(/usr/bin/codesign --verify --strict --verbose=2 "${bundle_path}" 2>&1)"; then
        log "codesign verification failed for ${bundle_label}: ${verify_output}"
        return 1
    fi

    if ! details_output="$(/usr/bin/codesign -dv --verbose=2 "${bundle_path}" 2>&1)"; then
        log "Couldn't inspect ${bundle_label} signature details: ${details_output}"
        return 1
    fi

    if echo "${details_output}" | grep -q 'Signature=adhoc'; then
        log "${bundle_label} is still ad hoc signed, which Tahoe may refuse during HAL reload."
        return 1
    fi

    if ! echo "${details_output}" | grep -q 'TeamIdentifier=' || \
       echo "${details_output}" | grep -q 'TeamIdentifier=not set'; then
        log "${bundle_label} does not have a TeamIdentifier. Refusing unsafe HAL reload."
        return 1
    fi

    return 0
}

validate_installed_products_before_hal_restart() {
    if [[ "${BGM_ALLOW_UNVERIFIED_HAL_RESTART:-0}" == "1" ]]; then
        log "WARNING: Skipping HAL restart signature validation because BGM_ALLOW_UNVERIFIED_HAL_RESTART=1."
        return 0
    fi

    verify_installed_bundle_for_hal_restart "/Library/Audio/Plug-Ins/HAL/Background Music Device.driver" "Background Music Device.driver" || return 1
    verify_installed_bundle_for_hal_restart "${xpc_helper_install_path}/BGMXPCHelper.xpc" "BGMXPCHelper.xpc" || return 1
    verify_installed_bundle_for_hal_restart "${dest_volume_no_trailing_slash}/Applications/Background Music.app" "Background Music.app" || return 1
}

find_browser_audio_session_pids() {
    local owner_uid
    owner_uid="$(stat -f %u /dev/console 2>/dev/null || echo 0)"

    ps -axo pid=,uid=,args= | \
        awk -v owner_uid="${owner_uid}" '
            $2 == owner_uid {
                cmd = ""
                for (i = 3; i <= NF; ++i) {
                    cmd = cmd $i " "
                }

                if (cmd ~ /audio\.mojom\.AudioService/ ||
                    cmd ~ /com\.apple\.WebKit\.(GPU|WebContent)/) {
                    print $1
                }
            }
        ' | sort -u
}

declare -a SUSPENDED_BROWSER_AUDIO_PIDS=()

suspend_browser_audio_clients_before_hal_restart() {
    SUSPENDED_BROWSER_AUDIO_PIDS=()

    while IFS= read -r pid; do
        [[ -n "${pid}" ]] && SUSPENDED_BROWSER_AUDIO_PIDS+=("${pid}")
    done < <(find_browser_audio_session_pids)

    if [[ ${#SUSPENDED_BROWSER_AUDIO_PIDS[@]} -eq 0 ]]; then
        return 0
    fi

    log "Temporarily suspending browser audio clients before restarting coreaudiod: ${SUSPENDED_BROWSER_AUDIO_PIDS[*]}"
    kill -STOP "${SUSPENDED_BROWSER_AUDIO_PIDS[@]}" 2>/dev/null || true
    sleep 1
}

resume_suspended_browser_audio_clients() {
    if [[ ${#SUSPENDED_BROWSER_AUDIO_PIDS[@]} -eq 0 ]]; then
        return 0
    fi

    log "Resuming browser audio clients after the guarded coreaudiod restart."
    kill -CONT "${SUSPENDED_BROWSER_AUDIO_PIDS[@]}" 2>/dev/null || true
    SUSPENDED_BROWSER_AUDIO_PIDS=()
}

wait_for_coreaudiod_restart() {
    local old_pid="$1"
    local old_start_time="$2"
    local timeout_seconds="${3:-20}"
    local waited=0

    while [[ ${waited} -lt ${timeout_seconds} ]]; do
        local new_pid
        new_pid="$(pgrep -x coreaudiod | head -n1 || true)"

        if [[ -n "${new_pid}" ]]; then
            local new_start_time
            new_start_time="$(ps -o lstart= -p "${new_pid}" 2>/dev/null | sed 's/^ *//')"

            if [[ -z "${old_pid}" ]] || [[ "${new_pid}" != "${old_pid}" ]] || \
               [[ -n "${old_start_time}" && "${new_start_time}" != "${old_start_time}" ]]; then
                sleep 2
                return 0
            fi
        fi

        sleep 1
        waited=$((waited + 1))
    done

    return 1
}

restart_coreaudiod_guarded() {
    local old_pid
    old_pid="$(pgrep -x coreaudiod | head -n1 || true)"
    local old_start_time=""
    if [[ -n "${old_pid}" ]]; then
        old_start_time="$(ps -o lstart= -p "${old_pid}" 2>/dev/null | sed 's/^ *//')"
    fi

    suspend_browser_audio_clients_before_hal_restart
    trap 'resume_suspended_browser_audio_clients' RETURN

    (sudo launchctl kickstart -k system/com.apple.audio.coreaudiod &>/dev/null || \
        sudo launchctl kill SIGTERM system/com.apple.audio.coreaudiod &>/dev/null || \
        sudo launchctl kill TERM system/com.apple.audio.coreaudiod &>/dev/null || \
        sudo launchctl kill 15 system/com.apple.audio.coreaudiod &>/dev/null || \
        sudo launchctl kill -15 system/com.apple.audio.coreaudiod &>/dev/null || \
        sudo killall coreaudiod &>/dev/null || \
        (sudo launchctl unload "$coreaudiod_plist" &>/dev/null && \
            sudo launchctl load "$coreaudiod_plist" &>/dev/null)) || return 1

    wait_for_coreaudiod_restart "${old_pid}" "${old_start_time}" 20 || return 1

    resume_suspended_browser_audio_clients
    trap - RETURN
}

coreaudiod_plist="/System/Library/LaunchDaemons/com.apple.audio.coreaudiod.plist"
dest_volume="$3"
dest_volume_no_trailing_slash="${dest_volume/%\//}"
xpc_helper_install_path="$(bash safe_install_dir.sh -y)"

log "Installing BGMXPCHelper to $xpc_helper_install_path"

# Delete the old version first, if any, to work around a macOS bug where a codesigned Mach-O binary
# can fail verification if you overwrite the old version instead of replacing it. Apparently the
# kernel caches some codesigning info that doesn't get updated. See
# <https://forums.developer.apple.com/thread/126187>.
rm -rf "${xpc_helper_install_path}/BGMXPCHelper.xpc"
cp -Rf "BGMXPCHelper.xpc" "$xpc_helper_install_path"

# Remove the quarantine attribute from BGMXPCHelper. Since macOS 14.0, this is needed to keep
# Gatekeeper from trying to show a permission prompt when BGMXPCHelper is launched. That would
# fail because BGMXPCHelper is a daemon, and Gatekeeper would prevent it from launching.
#
# TODO: I think, ideally, we should use the new Service Management API to install BGMXPCHelper,
#       which would probably avoid this issue. See
#       <https://developer.apple.com/documentation/servicemanagement/updating_your_app_package_installer_to_use_the_new_service_management_api>.
log "Removing com.apple.quarantine attribute from BGMXPCHelper."
sudo xattr -dr com.apple.quarantine "${xpc_helper_install_path}/BGMXPCHelper.xpc" 2>&1 | log_output
ls -lah@ "${xpc_helper_install_path}/BGMXPCHelper.xpc" 2>&1 | log_output

# TODO: Fail the install and show an error message if this fails.
bash "post_install.sh" "$xpc_helper_install_path" "BGMXPCHelper.xpc/Contents/MacOS/BGMXPCHelper" "."

# TODO: Verify the installed files, their permissions, the _BGMXPCHelper user/group, etc.

if ! validate_installed_products_before_hal_restart; then
    log "Refusing to restart coreaudiod because the installed HAL artifacts failed signature validation."
    exit 1
fi

log "Restarting coreaudiod with browser audio clients quiesced to avoid Tahoe HAL deadlocks."
if ! restart_coreaudiod_guarded; then
    log "Guarded coreaudiod restart failed. Suspended browser audio clients were resumed."
    exit 1
fi

# Wait until coreaudiod has restarted and BGMDevice is ready to use.
retries=5
while [[ $retries -gt 0 ]]; do
    device_list=""
    input_device_list=""

    if capture_command_with_timeout 8 device_list system_profiler SPAudioDataType && \
       [[ "$device_list" =~ "Background Music" ]]; then
        log "Background Music device is installed and available."
        retries=0
    elif capture_command_with_timeout 8 input_device_list ./ListInputDevices && \
        [[ "$input_device_list" =~ "BGMDevice" ]]; then
        log "Background Music device is installed and available."
        retries=0
    else
        retries=$((retries - 1))
        if [[ $retries -gt 0 ]]; then
            log "Background Music device not found. Trying again in 3 seconds..."
            sleep 3
        else
            # TODO: We might be able to use <installation-check> to show error messages in the
            #       installer GUI instead of just logging them. See
            #       <https://developer.apple.com/library/archive/documentation/DeveloperTools/Reference/DistributionDefinitionRef/Chapters/Distribution_XML_Ref.html#//apple_ref/doc/uid/TP40005370-CH100-SW12>.
            log "Background Music device not found. Installation failed."
            log "Audio devices:"
            log "$device_list"
            log "Audio input devices:"
            log "$input_device_list"
            exit 1
        fi
    fi
done

# Launch BGMApp.
#
# Allow a few retries because it seems to sometimes fail intermittently in CI builds. For example,
# <https://travis-ci.org/github/kyleneideck/BackgroundMusic/jobs/764126689>. From the error
# messages, it looks like BGMApp isn't always registered with Launch Services by this point.
logged_in_user_id="$(id -u "${USER}")"
did_open_bgmapp=false
retries=5

while [[ $did_open_bgmapp != "true" ]] && [[ $retries -gt 0 ]]; do
    retries=$((retries - 1))

    # Try opening BGMApp using its bundle ID first so the installer can declare it as "relocatable".
    # That way, if the user moves BGMApp and then installs a newer version of Background Music, the
    # installer will try to find the old version of BGMApp and put the new one in the same place.
    #
    # Use launchctl to make sure we don't run BGMApp as the installer user.
    #
    # TODO: If they have multiple copies of BGMApp, this might open one of the old ones.
    log "Opening Background Music.app by bundle ID"
    if launchctl asuser "${logged_in_user_id}" \
            open -b com.bearisdriving.BGM.App; then
        did_open_bgmapp=true
    fi

    if [[ $did_open_bgmapp != "true" ]]; then
        dest_volume_no_trailing_slash="${dest_volume/%\//}"
        log "Opening ${dest_volume_no_trailing_slash}/Applications/Background Music.app"
        if launchctl asuser "${logged_in_user_id}" \
                open "${dest_volume_no_trailing_slash}/Applications/Background Music.app"; then
            did_open_bgmapp=true
        fi
    fi

    if [[ $did_open_bgmapp != "true" ]]; then
        log "Opening Background Music.app using AppleScript"
        if osascript -e 'tell application "Background Music" to activate'; then
            did_open_bgmapp=true
        fi
    fi

    if [[ $did_open_bgmapp != "true" ]]; then
        log "Failed to open Background Music.app. Will retry shortly."
        sleep 1
    fi
done

# If we couldn't open BGMApp, it probably didn't install properly, but it's not certain, so don't
# fail the install.
if [[ $did_open_bgmapp != "true" ]]; then
    log "Failed to open Background Music.app. Giving up."
fi

# The installer plays a sound when it finishes, so give BGMApp a second to launch.
sleep 1

exit 0
